# comics/supervisor.conf
[program:gunicorn]
command=gunicorn comics.wsgi -b 0.0.0.0:8000
directory=/opt/django
user=nobody
autostart=true
autorestart=true
redirect_stderr=true